# supervisor/tk-dashboard-loos.conf
[program:loos-board-ui]
command=/opt/tk-dashboard/ui-apps/loos-ui-app.py
directory=/opt/tk-dashboard/ui-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0

[program:loos-import-io]
command=/opt/tk-dashboard/io-apps/loos-import-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=board-io
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"

[program:loos-meters-io]
command=/opt/tk-dashboard/io-apps/loos-meters-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=board-io
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"

[program:loos-screen-schedule]
command=/opt/tk-dashboard/io-apps/loos-screen-schedule-app.py /dev/serial/by-id/usb-FTDI_FT232R_USB_UART_AB0NDJEG-if00-port0
directory=/opt/tk-dashboard/io-apps/
user=board-io
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"
