# supervisor/tk-dashboard-wam.conf
[program:wam-board-ui]
command=/opt/tk-dashboard/ui-apps/wam-ui-app.py
directory=/opt/tk-dashboard/ui-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0

[program:wam-import-io]
command=/opt/tk-dashboard/io-apps/wam-import-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"
