# supervisor/tk-dashboard-loos.conf
[program:loos-board-ui]
command=/opt/tk-dashboard/ui-apps/loos-ui-app.py
directory=/opt/tk-dashboard/ui-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0

[program:loos-export-io]
command=/opt/tk-dashboard/io-apps/loos-export-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", python=/opt/tk-dashboard/io-apps/venv/bin/python

[program:loos-import-io]
command=/opt/tk-dashboard/io-apps/loos-import-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", python=/opt/tk-dashboard/io-apps/venv/bin/python

[program:loos-meters-io]
command=/opt/tk-dashboard/io-apps/loos-meters-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", python=/opt/tk-dashboard/io-apps/venv/bin/python
