# supervisor/tk-dashboard-mag.conf
[program:mag-board-ui]
command=/opt/tk-dashboard/ui-apps/mag-ui-app.py
directory=/opt/tk-dashboard/ui-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0

[program:mag-import-io]
command=/opt/tk-dashboard/io-apps/mag-import-io-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"

[program:wam-screen-schedule]
command=/opt/tk-dashboard/io-apps/mag-screen-schedule-app.py
directory=/opt/tk-dashboard/io-apps/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"
