# supervisor/tk-dashboard-loos.conf
[program:loos-board-ui]
command=/opt/tk-dashboard/loos-board-ui/app.py
directory=/opt/tk-dashboard/loos-board-ui/
user=pi
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1", DISPLAY=:0.0
